# bin/Debug/netcoreapp2.2/linux-x64-rhel8/golink.service
[Unit]
Description=golink

[Service]
Type=Simple
ExecStart=/var/www/golink/golink --urls "http://0.0.0.0:80"
WorkingDirectory=/var/www/golink/publish
Environment=ASPNETCORE_ENVIRONMENT=Production
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=golink

[Install]
WantedBy=multi-user.target
