# bin/Debug/netcoreapp2.2/linux-x64/golink.service
[Unit]
Description=golink

[Service]
Type=notify
ExecStart=/var/www/golink/golink
User=svcgolink
WorkingDirectory=/var/www/golink/

[Install]
WantedBy=multi-user.target